# ops/bootstrap/macos/launchd/com.jcn.controlplane.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jcn.controlplane</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__BRINGUP_SCRIPT__</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__STDOUT_LOG__</string>
  <key>StandardErrorPath</key>
  <string>__STDERR_LOG__</string>
  <key>WorkingDirectory</key>
  <string>__WORKING_DIRECTORY__</string>
</dict>
</plist>
