# ops/bootstrap/macos/launchd/com.jcn.controlplane.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jcn.controlplane</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/justin/Documents/Justyn Clark Network/REPOS/homelab-control-plane/ops/bootstrap/macos/bringup.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>/Users/justin/Documents/Justyn Clark Network/REPOS/homelab-control-plane/receipts/launchd/bringup.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/justin/Documents/Justyn Clark Network/REPOS/homelab-control-plane/receipts/launchd/bringup.stderr.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/justin/Documents/Justyn Clark Network/REPOS/homelab-control-plane</string>
</dict>
</plist>

